# deploy_tools/gunicorn-upstart.template.conf
description "Gunicorn server for SITENAME"

start on net-device-up
stop on shutdown

respawn

setuid jon
chdir /home/jon/sites/SITENAME
exec export PYTHONPATH=${HOME}/sites/SITENAME/SITENAME

exec /home/jon/.virtualenvs/SITENAME/bin/gunicorn \
--bind unix:/tmp/SITENAME.socket
config.wsgi:application
